# etc/services.d/php-fpm81/run
#!/command/with-contenv sh
php-fpm81 -F -O -c /etc/php81/php-fpm.conf
