# etc/services.d/php-fpm82/run
#!/command/with-contenv sh
php-fpm82 -F -O -c /etc/php82/php-fpm.conf
