# etc/services.d/php-fpm8/run
#!/command/execlineb -P
php-fpm8 -F -R -c /etc/php8/php-fpm.conf
