# etc/services.d/php-fpm8/run
#!/command/with-contenv sh
php-fpm8 -F -O -c /etc/php8/php-fpm.conf
