# etc/services.d/php-fpm/run
#!/command/with-contenv sh
$PHP_FPM_BIN -F -O -c /etc/php/$PHP_VERSION/fpm -y /etc/php/${PHP_VERSION}/fpm/php-fpm.conf
