# etc/services.d/php-fpm5/run
#!/command/with-contenv sh
php-fpm5 -F -O -c /etc/php5/php-fpm.conf
