# etc/services.d/php-fpm5/run
#!/command/execlineb -P
php-fpm5 -F -R -c /etc/php5/php-fpm.conf
