# etc/services.d/php-fpm7/run
#!/command/execlineb -P
php-fpm7 -F -R -c /etc/php7/php-fpm.conf
