# etc/services.d/php-fpm7/run
#!/command/with-contenv sh
php-fpm7 -F -O -c /etc/php7/php-fpm.conf
